# net-analyzer/gvmd/files/gvmd.service
[Unit]
Description=Greenbone Vulnerability Manager
After=network.target
After=openvas-scanner.service
Wants=openvas-scanner.service
Before=gsad.service

[Service]
Type=forking
User=gvm
Group=gvm
EnvironmentFile=-/etc/gvm/sysconfig/gvmd-daemon.conf
ExecStart=/usr/bin/gvmd $GVMD_OPTIONS $GVMD_PORT $GVMD_LISTEN_ADDRESS $GVMD_SCANNER_HOST $GVMD_GNUTLS_PRIORITIES 
Restart=on-failure
RestartSec=12

[Install]
WantedBy=multi-user.target
