# net-analyzer/gvmd/files/gvmd.service
[Unit]
Description=OpenVAS Manager
After=network.target
After=openvas-scanner.service
Wants=openvas-scanner.service
Before=gsa.service

[Service]
Type=forking
EnvironmentFile=-/etc/gvm/sysconfig/gvmd-daemon.conf
ExecStart=/usr/sbin/gvmd $GVMD_OPTIONS $GVMD_PORT $GVMD_LISTEN_ADDRESS $GVMD_SCANNER_HOST $GVMD_GNUTLS_PRIORITIES 
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
User=root
Group=root
TimeoutSec=1200

[Install]
WantedBy=multi-user.target
